# supervisord.conf
[supervisord]
nodaemon=true
user=root

# 启动 VNC Server
[program:vnc]
command=/bin/bash -c "rm -rf /tmp/.X1-lock /tmp/.X11-unix/X1 && vncserver :1 -geometry 1920x1080 -depth 24 -localhost no -rfbauth /root/.vnc/passwd"
autostart=true
autorestart=true
stderr_logfile=/var/log/vnc.err.log
stdout_logfile=/var/log/vnc.out.log

# 启动 noVNC (Web 代理)
[program:novnc]
command=websockify -D --web=/usr/share/novnc/ 6080 localhost:5901
autostart=true
autorestart=true
stderr_logfile=/var/log/novnc.err.log
stdout_logfile=/var/log/novnc.out.log

# 启动 Telegram
[program:telegram]
# 使用 dbus-run-session 确保 Telegram 能正确连接会话总线
command=dbus-run-session telegram-desktop
autostart=true
autorestart=true
user=root
environment=HOME="/root",USER="root"
stderr_logfile=/var/log/telegram.err.log
stdout_logfile=/var/log/telegram.out.log
